# deploy/systemd/signalengine-signals.timer
[Unit]
Description=Run SignalEngine retrain + signals after the nightly ingest

[Timer]
OnCalendar=Mon..Sat 23:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
